# bind/bind-update-zones.cron.d
@weekly root /usr/local/bin/bind-update-zones && kill -HUP $(pidof named)
